# propart.service
[Unit]
Description=ProPart Real Estate Hub - Gunicorn
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=propart
WorkingDirectory=/home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub
Environment="PATH=/home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub/venv/bin"
EnvironmentFile=/home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub/.env
ExecStart=/home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub/venv/bin/gunicorn \
          --config /home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub/gunicorn_config.py \
          --bind 127.0.0.1:8000 \
          run:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pro-part-agent/htdocs/agent.pro-part.online/propart-real-estate-hub /var/log/propart /var/run/propart

[Install]
WantedBy=multi-user.target

